# builds/jacktheflipper/kiosk/ttyd-claude.service
[Unit]
Description=Claude CLI exposed as web terminal (ttyd) on JackTheFlipper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jacktheflipper
WorkingDirectory=/home/jacktheflipper
Environment=HOME=/home/jacktheflipper TERM=xterm-256color
ExecStart=/usr/local/bin/ttyd -W -p 7777 -t titleFixed=Claude /usr/local/bin/claude
Restart=on-failure
RestartSec=4

[Install]
WantedBy=multi-user.target
