# builds/node1-stack/sinsera-backup.service
[Unit]
Description=Sinsera self-hosted stack backup (Supabase DB + Gitea + config)
[Service]
Type=oneshot
ExecStart=/opt/sinsera-stack/backup.sh
