# builds/pi-zero2-kiosk/pi-zero2-kiosk-install.service
[Unit]
Description=Sinsera Kiosk phase-2 install (Chromium + X11 + autostart)
Documentation=https://sinsera.co/ark/
After=network-online.target time-sync.target
Wants=network-online.target
ConditionPathExists=/opt/pi-zero2-kiosk/install.sh

[Service]
Type=oneshot
ExecStart=/opt/pi-zero2-kiosk/install.sh
RemainAfterExit=yes
TimeoutSec=1800
StandardOutput=append:/var/log/pi-zero2-kiosk-install.log
StandardError=append:/var/log/pi-zero2-kiosk-install.log

[Install]
WantedBy=multi-user.target
