# apps/the-comb/the-comb.service
[Unit]
Description=The Comb launcher app (zero-dep Node server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=brocoli
WorkingDirectory=/opt/the-comb
Environment=NODE_ENV=production PORT=8080 HOST=0.0.0.0
ExecStart=/usr/bin/node /opt/the-comb/server.mjs
Restart=on-failure
RestartSec=4

[Install]
WantedBy=multi-user.target
