# builds/the-hive/the-hive-install.service
[Unit]
Description=The Hive phase-2 install (Claude CLI + Tor + X + Chromium + The Comb)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/the-hive/install.sh

[Service]
Type=oneshot
ExecStart=/opt/the-hive/install.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
